# backend/deploy/supervisor.conf
; Supervisor config for license-manager-v1 gunicorn process.
;
; Install:
;   sudo cp /home/django/license-manager/backend/deploy/supervisor.conf \
;            /etc/supervisor/conf.d/license-manager-v1.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;
; Legacy supervisor program is "license-manager" — this is a distinct program.
; Both run in parallel during the parallel-run phase.
;
; At cutover:
;   1. Update gunicorn.conf.py bind to "127.0.0.1:8000" (or set GUNICORN_BIND).
;   2. Stop legacy: sudo supervisorctl stop license-manager license-manager-celery
;   3. Reload nginx to point / at backendv1 (remove the /api/v1/ split routing).

[program:license-manager-v1]
command=/home/django/license-manager/venv/bin/gunicorn config.wsgi:application -c /home/django/license-manager/backend/deploy/gunicorn.conf.py
directory=/home/django/license-manager/backend
user=django
autostart=true
autorestart=true
; Allow up to 30 s for gunicorn to finish in-flight requests before SIGKILL.
stopwaitsecs=30
stopsignal=TERM
stderr_logfile=/var/log/license-manager-v1/supervisor-error.log
stdout_logfile=/var/log/license-manager-v1/supervisor-out.log
; Ensure /var/log/license-manager-v1/ exists before starting:
;   sudo mkdir -p /var/log/license-manager-v1 && sudo chown django: /var/log/license-manager-v1
environment=DJANGO_SETTINGS_MODULE="config.settings.prod"
